# deploy/systemd/user/ai-lab-dashboard-smoke.timer
[Unit]
Description=Run AI Lab Dashboard smoke test daily

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
